# systemd/hermes-satellite.service
[Unit]
Description=hermes-satellite voice front-end for Hermes
# Wait for sound devices (ALSA codec driver / device-tree overlay) to be ready.
After=sound.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# Dedicated unprivileged user; must be in the spi, audio and gpio groups so it
# can drive the APA102 LEDs (SPI), the microphone/speaker (ALSA) and the HAT
# button (GPIO). Create with:
#   sudo useradd -r -G spi,audio,gpio hermes
User=hermes
Group=hermes
SupplementaryGroups=spi audio gpio

# Adjust paths to your install. The venv is created at /opt/hermes-satellite.
WorkingDirectory=/opt/hermes-satellite
# Re-apply the WM8960 mixer state on every start — defends against stale
# alsactl snapshots silently disabling the capture path (see the pi4 hardware
# guide's troubleshooting section). Uncomment on ReSpeaker 2-Mic v1 setups:
#ExecStartPre=/opt/hermes-satellite/scripts/wm8960-mixer.sh
ExecStart=/opt/hermes-satellite/.venv/bin/python -m hermes_satellite \
    --config /etc/hermes-satellite/config.yaml
# Secrets can be provided out-of-band instead of in config.yaml:
#   EnvironmentFile=/etc/hermes-satellite/secrets.env
# containing HERMES_API_KEY=..., HERMES_SESSION_KEY=..., PORCUPINE_ACCESS_KEY=...
EnvironmentFile=-/etc/hermes-satellite/secrets.env

# The service user has no home directory, so point XDG caches (Moonshine's
# ~80 MB STT model lives in $XDG_CACHE_HOME/moonshine_voice) at the app's
# data directory. Pre-seed it to avoid a download at first service start:
#   sudo -u hermes XDG_CACHE_HOME=/var/lib/hermes-satellite/cache \
#     /opt/hermes-satellite/.venv/bin/python -c \
#     "import moonshine_voice as mv; mv.get_model_for_language('en', mv.ModelArch.BASE)"
Environment=XDG_CACHE_HOME=/var/lib/hermes-satellite/cache

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
